# src/main/bash/crontab

SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
HOME=/home/pi

* 3 * * * iot-upgrade > $HOME/upgrade.log 2>&1
* 4 * * * sudo reboot

# needs newline at the end
